# runit/mqttwrk/run
#!/bin/sh

exec /usr/share/bytebeam/mqttwrk/target/release/mqttwrk bench -n 1000000 -a 1000 -b 1 -r 1 -h beamd
